# services.d/sshd/run
#!/usr/bin/with-contenv /bin/bash

exec sshd -D

